# systemd/spot-dashboard.service.d/override.conf
[Service]
EnvironmentFile=/home/ubuntu/Spot-Ai-Agent/.env
WorkingDirectory=/home/ubuntu/Spot-Ai-Agent
ReadWritePaths=/home/ubuntu/spot_data /home/ubuntu/spot_data/logs /home/ubuntu/spot_data/trades
ProtectHome=off
StandardOutput=journal
StandardError=inherit
